# xyz.openbmc_project.MonitoringService.service.in
[Unit]
Description=Monitoring Service
After=dbus.service

[Service]
ExecStart=@bindir@/monitoring-service
Restart=always
Type=dbus
BusName=xyz.openbmc_project.MonitoringService

[Install]
WantedBy=multi-user.target
