# xyz.openbmc_project.Telemetry.service.in
[Unit]
Description=Telemetry@wants@@requires@@before@
After=dbus.service@after@

[Service]
ExecStart=@bindir@/telemetry
Restart=always
Type=dbus
BusName=xyz.openbmc_project.Telemetry

[Install]
WantedBy=multi-user.target
